# datasette.cron
0 6 * * 2 root /usr/bin/systemctl restart datasette.service