# datasette.cron
0 6 * * 2 /usr/bin/systemctl restart datasette.service