# datasette.cron
0 6 * * * /usr/bin/systemctl restart datasette.service